# deploy/jorasite.service
[Unit]
Description=HEAVY / jorasite — Node app (Express + Meta CAPI + pCloud proxy)
Documentation=file:/opt/jorasite/README.md
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=jorasite
Group=jorasite
WorkingDirectory=/opt/jorasite/app
EnvironmentFile=/opt/jorasite/.env
ExecStart=/usr/bin/node /opt/jorasite/app/server.js
Restart=on-failure
RestartSec=3

# ── Logging → journald (view: journalctl -u jorasite) ──
StandardOutput=journal
StandardError=journal
SyslogIdentifier=jorasite

# ── Hardening: minimal privileges for the app process ──
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
# Everything is read-only (ProtectSystem=strict) EXCEPT the data dir, where
# the app persists orders + issued tickets (JSON store).
ReadWritePaths=/opt/jorasite/data
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_INET AF_INET6
RestrictNamespaces=true
LockPersonality=true
MemoryDenyWriteExecute=false
RestrictSUIDSGID=true

[Install]
WantedBy=multi-user.target
